# roles/setupChatBotService/tasks/templates/chatbot.timer.j2
[Unit]
Description="Run script after 15 seconds after boot"

[Timer]
OnBootSec=15s

[Install]
WantedBy=default.target
