# roles/setupABS/templates/chatbot.service
[Unit]
Description=Run script at startup

[Service]
Type=oneshot
ExecStart=/etc/init.d/chatbot
TimeoutStartSec=0
